# .config/systemd/user/check_update.service
[Unit]
Description=check update 

[Service]
Type=simple
ExecStart=%h/bin/check_update

[Install]
WantedBy=multi-user.target