# .config/systemd/user/check_update.service
[Unit]
Description=check update

[Service]
Type=oneshot
ExecStart=%h/bin/check_update

